# deploy/platform_server/asdun-platform.service
[Unit]
Description=ASDUN Platform Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=asdun
Group=asdun
WorkingDirectory=/opt/asdun_pi/platform_server
EnvironmentFile=/etc/asdun/platform.env
Environment=ASDUN_PLATFORM_DB=/var/lib/asdun-platform/asdun_platform.sqlite
ExecStart=/opt/asdun_pi/platform_server/.venv/bin/python -m uvicorn app:app --host 127.0.0.1 --port 9000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
